# supervisord.conf
[supervisord]
nodaemon=true

[program:netatalk]
command=/usr/local/sbin/netatalk -d